# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/byobu
StandardInput=tty
StandardOutput=tty
User=st_lim
Group=st_lim
