# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=+/usr/bin/sudo -u st_lim -H -i /usr/bin/byobu
StandardInput=tty
StandardOutput=tty
User=root
